# init/wye.service.in

[Unit]
Description=wye streaming service
Documentation=http://wye.sourceforge.net/

[Service]
PIDFile=/var/run/wye.pid
ExecStart=@bindir@/wye-service

[Install]
WantedBy=multi-user.target

